# .config/autostart/syncthing.desktop
[Desktop Entry]
Type=Application
Name=Syncthing
Exec=syncthing --no-browser
