# .config/autostart/wpchg_random.desktop
[Desktop Entry]
Comment=
Exec=/media/stuff/vadim/bin/wpchg_random
Hidden=false
Name=/media/stuff/vadim/bin/wpchg_random
Type=Application
